# roles/manageiq.postgresql/templates/etc/profile.d/scl.sh.j2
#!/bin/bash
source /opt/rh/{{ scl }}/enable
export X_SCLS="`scl enable {{ scl }} 'echo $X_SCLS'`"
